# etc/systemd/system/paccache.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/paccache -rk1 ; /usr/bin/paccache -ruk0
